# scripts/deploy/legacy/systemd/vodeneevbet-parser.service
[Unit]
Description=VodeneevBet Parser Service
After=network.target

[Service]
Type=simple
User=vodeneevm
WorkingDirectory=REMOTE_DIR/internal/parser
ExecStart=REMOTE_DIR/internal/parser/parser -config REMOTE_DIR/configs/local.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vodeneevbet-parser

# Environment variables
Environment="GOPATH=/home/vodeneevm/go"
Environment="PATH=/usr/local/go/bin:/home/vodeneevm/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

