# scripts/deploy/systemd/vodeneevbet-api.service
[Unit]
Description=VodeneevBet API Service
After=network.target

[Service]
Type=simple
User=vodeneevm
WorkingDirectory=REMOTE_DIR/internal/api
ExecStart=REMOTE_DIR/internal/api/api
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vodeneevbet-api

# Environment variables
Environment="GOPATH=/home/vodeneevm/go"
Environment="PATH=/usr/local/go/bin:/home/vodeneevm/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
